# systemd/water-controller-api.service
[Unit]
Description=Water Treatment Controller Web API
Documentation=https://github.com/mwilco03/Water-Controller
After=network.target water-controller.service
BindsTo=water-controller.service

[Service]
Type=simple
User=wtc
Group=wtc
WorkingDirectory=/opt/water-controller/web/api
ExecStart=/opt/water-controller/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8080 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-controller-api

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/water-controller /var/log/water-controller /dev/shm

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/water-controller/environment

[Install]
WantedBy=multi-user.target
