# systemd/water-controller-api.service
[Unit]
Description=Water Treatment Controller API (FastAPI)
Documentation=https://github.com/mwilco03/Water-Controller
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=water-controller
Group=water-controller
WorkingDirectory=/opt/water-controller/app
ExecStart=/opt/water-controller/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1 --proxy-headers
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-controller-api

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/water-controller /var/log/water-controller /run/water-controller

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=CONFIG_PATH=/etc/water-controller/config.yaml
Environment=DATA_DIR=/var/lib/water-controller
Environment=LOG_DIR=/var/log/water-controller
EnvironmentFile=-/etc/water-controller/environment

# Resource limits
MemoryMax=512M
MemoryHigh=400M
CPUQuota=200%
TasksMax=256
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
