# systemd/water-controller-ui.service
[Unit]
Description=Water Treatment Controller Web UI
Documentation=https://github.com/mwilco03/Water-Controller
After=network.target water-controller-api.service
# Prefer API service but don't hard-fail if unavailable
Wants=water-controller-api.service

[Service]
Type=simple
User=wtc
Group=wtc
WorkingDirectory=/opt/water-controller/web/ui

# Pre-start: log build status but don't block startup
# State detection is handled by health endpoint and UI degraded banner
ExecStartPre=-/bin/bash -c '/opt/water-controller/scripts/lib/ui_build_check.sh --quick 2>&1 || echo "WARNING: UI build validation failed - system may be degraded"'

ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-controller-ui

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/water-controller
ReadOnlyPaths=/opt/water-controller

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=WTC_UI_PATH=/opt/water-controller/web/ui
EnvironmentFile=-/etc/water-controller/environment

[Install]
WantedBy=multi-user.target
