# systemd/water-controller-ui.service
[Unit]
Description=Water Treatment Controller Web UI (Next.js)
Documentation=https://github.com/mwilco03/Water-Controller
After=network.target water-controller-api.service
Wants=water-controller-api.service

# Readiness gate: Ensure UI has been built
# Anti-pattern addressed: "UI treated as optional"
# The UI is critical - operators cannot function without it
ConditionPathExists=/opt/water-controller/web/ui/.next/build-manifest.json
ConditionPathIsDirectory=/opt/water-controller/web/ui/.next/static

[Service]
Type=simple
User=wtc
Group=wtc
WorkingDirectory=/opt/water-controller/web/ui

# Startup validation: Verify UI build is complete
ExecStartPre=/bin/bash -c 'test -f /opt/water-controller/web/ui/.next/build-manifest.json || (echo "ERROR: UI not built - run npm run build" && exit 1)'
ExecStartPre=/bin/bash -c 'test -d /opt/water-controller/web/ui/.next/static || (echo "ERROR: UI static assets missing" && exit 1)'
ExecStartPre=/bin/bash -c 'test -d /opt/water-controller/web/ui/.next/server || (echo "ERROR: UI server build missing" && exit 1)'

# Main service - Port from WTC_UI_PORT environment variable
ExecStart=/bin/bash -c '/usr/bin/node /opt/water-controller/web/ui/node_modules/.bin/next start -p ${WTC_UI_PORT:-8080} -H 0.0.0.0'
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-controller-ui

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/water-controller /var/log/water-controller

# Environment - Ports from centralized config
# See config/ports.env for the single source of truth
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=NEXT_TELEMETRY_DISABLED=1
# Load centralized port configuration first
EnvironmentFile=-/opt/water-controller/config/ports.env
EnvironmentFile=-/etc/water-controller/environment

# Resource limits
MemoryMax=256M
MemoryHigh=200M
TasksMax=64

[Install]
WantedBy=multi-user.target
