# systemd/water-controller-hmi.service
[Unit]
Description=Water Treatment Controller HMI
Documentation=https://github.com/mwilco03/Water-Controller
After=network.target

[Service]
Type=simple
User=water-controller
Group=water-controller
WorkingDirectory=/opt/water-controller
ExecStart=/usr/local/bin/water-controller --web-only
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-controller-hmi

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/water-controller /var/log/water-controller

# Environment
Environment=DATA_DIR=/var/lib/water-controller
Environment=WEB_PORT=8080
Environment=PYTHONUNBUFFERED=1
Environment=NODE_ENV=production
EnvironmentFile=-/etc/water-controller/environment

[Install]
WantedBy=multi-user.target
